# kafka/supervisor/kafdrop.conf
[program:kafdrop]
command=%(ENV_JAVA_HOME)s/bin/java -jar /usr/local/bin/kafdrop.jar --zookeeper.connect=%(ENV_ZK_HOSTS)s --server.port=%(ENV_LISTEN)s
autostart=true
autorestart=true